# systemd/sentinel.service
[Unit]
Description=VPS Sentinel Security Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sentinel-monitor.sh
StandardOutput=journal
StandardError=journal
